# crontab
50 5,7,11,18 * * * cd ~/FutuOpenDClient/WeatherForecast; ./start.sh > start.txt 2>&1