# crontab
50 5,7,11,18 * * * cd ~/FutuOpenDClient/WeatherForecast; export PUPPETEER_EXECUTABLE_PATH=/opt/google/chrome/chrome; npm start