# weather/crontab
2 8,12,19 * * * cd ~/web-environment/map/weather; ./start.sh > start.log 2>&1
